# no-android.conf
start on startup

pre-start script
        if [ ! -f /var/lib/lxc/android/config ]; then
                initctl emit android
        fi
end script

