# ubuntu-session.conf
author "Tony Espy <espy@canonical.com>"

start on mounted MOUNTPOINT=/run

setuid phablet
setgid android_input

script
        
        exec /usr/bin/ubuntu-session

end script
