# debian/ubuntu-touch-session.upstart
author "Tony Espy <espy@canonical.com>"

start on (started lxc-android-config
	     and started dbus)
stop on stopping dbus

setuid phablet
setgid android_input
oom score -1

respawn

exec /usr/bin/ubuntu-touch-session
