# debian/ubuntu-touch-session.upstart
author "Tony Espy <espy@canonical.com>"

start on (started udev
	     and started dbus)
stop on stopping dbus

setuid phablet
setgid android_input
oom score -10

respawn

exec /usr/bin/ubuntu-touch-session
