# debian/ubuntu-session.upstart
author "Tony Espy <espy@canonical.com>"

start on mounted MOUNTPOINT=/run
stop on stopping dbus

setuid phablet
setgid android_input

respawn

script
        
        exec /usr/bin/ubuntu-session

end script
