# upstart-session/click-install.conf
description "Click app installer"
author "Sergio Schvezov <sergio.schvezov@ubuntu.com>"

console log

start on started dbus

emits click-install-complete

exec preinstall-click
