# upstart-session/friends-service.conf
description "Friends Service"
author "Ricardo Mendoza <ricmm@ubuntu.com>"

start on started dbus
stop on session-end

respawn

exec friends-service -do
