# upstart-session/maliit-server.conf
description "Maliit OSK Server"
author "Ricardo Mendoza <ricmm@ubuntu.com>"

start on started dbus
stop on session-end

respawn

script
    sleep 1
    exec maliit-server
end script
