# upstart-session/ofono-setup.conf
description "oFono setup"
author "Michał Sawicz <michal.sawicz@canonical.com>"

start on started unity8

script
    sleep 12
    /usr/bin/ofono-setup
end script
