# upstart-session/ofono-setup.conf
description "Setup oFono when unity8 starts"
author "Oliver Grawert <ogra@canonical.com>"

start on started dbus
stop on desktop-end

exec /usr/bin/ofono-setup
