# upstart-session/ofono-setup.conf
description "Setup telepathy accounts when ofono and dbus start"
author "Oliver Grawert <ogra@canonical.com>"

start on started dbus and started pulseaudio
stop on desktop-end

exec ofono-setup
