# upstart-session/ofono-setup.conf
description "Setup oFono when unity8 starts"
author "Oliver Grawert <ogra@canonical.com>"

start on (:sys:started JOB=ofono) and (started dbus)
stop on desktop-end

exec /usr/bin/ofono-setup
