# etc/profile.d/10-upstart-phablet.sh
# Join the phablet Upstart session
export UPSTART_SESSION=$(initctl list-sessions | awk '{ print $NF; quit }')
