# etc/profile.d/qt_networkmanager_bearer.sh
# Allow using the proper network manager bearer
export QT_EXCLUDE_GENERIC_BEARER=1
