# etc/profile.d/ubuntu-touch-file-selector.sh
export QT_FILE_SELECTORS="ubuntu-touch"
